# systemd/mon_tcp/pcap-dnsreader.service
[Service]
Type=oneshot
ExecStart=/tmp/dnsreader /mnt/2/pcap/53.pcap /srv/http/queries.json
