# systemd/mon_tcp/pcap-dnsreader.path
[Unit]
Description=Watch for changes
ConditionUser=!root
[Path]
PathModified=/mnt/2/pcap/53.pcap
